# ubuntu-playground/upstart-neos/debian/neos.conf
# neos - start neos

description "start neos"

start on (startup and filesystem and started udev and stopped udevtrigger)
stop on runlevel [!2345]

task
exec /bin/startneos
